# supervisord.conf
[program:node]
command = node ./lib/app.js
directory = /home/dotcloud/current

